# extras/.config/systemd/user/power-monitor.service
[Service]
Environment=STARTUP_WAIT=30s
ExecStart=~/.local/bin/power-monitor.sh

[Install]
WantedBy=default.target
